# linux_configuration/scripts/single_use/features/systemd/keepass-cloud-sync.service
[Unit]
# Mirror the live KeePassXC vault into the dufs-served cloud folder.
# Triggered by keepass-cloud-sync.path on every change to Passwords.kdbx.
# Placeholders (__USER__ / __HOME__ / __CLOUD_ROOT__) are substituted by
# setup_dufs_cloud.sh at install time.
Description=Mirror KeePass vault into dufs cloud folder

[Service]
Type=oneshot
User=__USER__
# Only the vault file — not the whole ~/Keepass dir (keeps local .backup_* vaults
# off the public server). No --delete: the cloud folder may hold other files.
ExecStart=/usr/bin/rsync -a __HOME__/Keepass/Passwords.kdbx __CLOUD_ROOT__/Keepass/
