# phone_focus_mode/systemd/phone-auto-sync.service
[Unit]
Description=Phone focus mode auto sync
After=network.target

[Service]
Type=oneshot
ExecStart=%h/testsAndMisc/scripts/run_all/run_phone.sh auto
StandardOutput=journal
StandardError=journal
