# linux_configuration/scripts/periodic_background/hosts/guard/hosts-guard.path
[Unit]
Description=Watch /etc/hosts and trigger enforcement

[Path]
PathChanged=/etc/hosts
Unit=hosts-guard.service

[Install]
WantedBy=multi-user.target
