# linux_configuration/scripts/periodic_background/hosts/guard/nsswitch-guard.path
[Unit]
Description=Watch /etc/nsswitch.conf for tampering (hosts bypass protection)

[Path]
PathChanged=/etc/nsswitch.conf
Unit=nsswitch-guard.service

[Install]
WantedBy=multi-user.target
