# python_pkg/steam_backlog_enforcer/steam-backlog-enforcer.service
[Unit]
Description=Steam Backlog Enforcer
After=network-online.target graphical.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/kuhy/testsAndMisc
ExecStart=/usr/bin/python3 -m python_pkg.steam_backlog_enforcer.main enforce
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/kuhy/.local/lib/python3.14/site-packages
Environment=HOME=/home/kuhy
# Hardening: enforcer must not be easily killed.
OOMScoreAdjust=-900

[Install]
WantedBy=multi-user.target
