# linux_configuration/scripts/periodic_background/system-maintenance/systemd/auto-system-update.service
[Unit]
Description=Automatic System Update (pacman + yay AUR)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/bin/auto-system-update.sh
StandardOutput=journal
StandardError=journal

# System updates can take a while on slow connections
TimeoutStartSec=1800
TimeoutStopSec=30

Restart=no

[Install]
WantedBy=multi-user.target
